# meta-edtwardy/recipes-connectivity/bluez-iot-agent/bluez-iot-agent/bluez-iot-agent.service
[Unit]
Description=BlueZ pairing wizard for IoT devices
After=bluetooth.service network.target
Requires=bluetooth.service

[Service]
Type=dbus
BusName=org.bluez.iot-agent
ExecStartPre=/bin/sleep 7
ExecStart=/usr/bin/bluez-iot-agent

[Install]
WantedBy=multi-user.target
