# layers/meta-email-gadget/recipes-daemons/postfix/files/postfix-makeconf.service
[Unit]
Description=postfix configuration
Before=postfix.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/postfix-makeconf.sh

[Install]
RequiredBy=postfix.service
