# layers/meta-backup-gadget/recipes-core/backup-dataset/files/backup-controller.service
[Unit]
Description=Backup Controller

[Service]
ExecStart=/usr/bin/backup-controller
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
